# libubic-service-php5-fcgi-perl/debian/postinst
#!/bin/sh

if [ -x '/etc/init.d/php-fcgi' ]
then
	update-rc.d -f nginx remove
fi

ubic start php5-fcgi

