# libubic-service-php5-fcgi-perl/debian/postinst
#!/bin/sh

if [ -x '/etc/init.d/php5-fcgi' ]
then
	update-rc.d -f php5-fcgi remove
fi

invoke-rc.d php5-fcgi restart
