# libubic-service-nginx-perl/debian/postinst
#!/bin/sh

if [ -x '/etc/init.d/nginx' ]
then
	update-rc.d -f nginx remove
fi

invoke-rc.d nginx restart
