# libubic-service-nginx-perl/debian/prerm
#!/bin/sh

if [ -x "/etc/init.d/nginx" ]
then
	invoke-rc.d nginx stop
fi
