# libubic-service-haproxy-perl/debian/postinst
#!/bin/sh

if [ -x '/etc/init.d/haproxy' ]
then
	update-rc.d -f haproxy remove
fi

ubic start haproxy

